# deploy/supervisor/quantoffice.conf
# =============================================================================
# QuantOffice Supervisor 进程托管
# -----------------------------------------------------------------------------
# 拷贝位置：
#   /etc/supervisor/conf.d/quantoffice.conf
#
# 启用：
#   sudo cp deploy/supervisor/quantoffice.conf /etc/supervisor/conf.d/quantoffice.conf
#   sudo mkdir -p /var/log/supervisor /var/run
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start quantoffice
#   sudo supervisorctl status quantoffice
#
# 依赖：先跑过 ./deploy.sh install（生成 /workspace/.venv）
# 端口：监听 127.0.0.1:8765（只对本地，nginx 转发）
# 进程：单 worker（uvicorn）—— QuantOffice 内部用 async，1 worker 已足够
# =============================================================================

[program:quantoffice]
; 启动命令 —— 用 venv 里的 uvicorn
command=/workspace/.venv/bin/uvicorn quant_office.app:app
    --host 127.0.0.1
    --port 8765
    --log-level info

; 项目根目录（Settings._project_root() 读 CWD 时用得到）
directory=/workspace

; 运行用户（改为你服务器上的非 root 用户；root 也能跑但不推荐）
user=www-data

; 自启 & 崩溃自愈
autostart=true
autorestart=true
startretries=10
startsecs=5

; 整组杀进程（uvicorn 子进程一并收掉）
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

; 优先级
priority=100

; 环境变量：PYTHONUNBUFFERED 强制 flush，PATH 包含 venv
environment=
    PYTHONUNBUFFERED="1",
    PYTHONDONTWRITEBYTECODE="1",
    PATH="/workspace/.venv/bin:/usr/local/bin:/usr/bin:/bin"

; 日志（先确保 /var/log/supervisor 存在）
stdout_logfile=/var/log/supervisor/quantoffice.out.log
stderr_logfile=/var/log/supervisor/quantoffice.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=10
stdout_logfile_relevant=true
