# deploy/supervisor/quantoffice-frontend.conf
# =============================================================================
# QuantOffice 前端 — Vite dev 服务器（supervisor 托管）
# -----------------------------------------------------------------------------
# 适用场景：
#   - 改了前端代码想即时热重载（Vite HMR）
#   - 1H1G 机器跑不动 bun run build（dev 模式内存友好）
#   - 不想每次改前端都跑 ./deploy.sh install
#
# 拷贝位置：
#   /etc/supervisor/conf.d/quantoffice-frontend.conf
#
# 启用：
#   sudo mkdir -p /var/log/supervisor /opt/quantoffice/scripts
#   sudo cp deploy/scripts/start-frontend-dev.sh /opt/quantoffice/scripts/
#   sudo chmod +x /opt/quantoffice/scripts/start-frontend-dev.sh
#   sudo cp deploy/supervisor/quantoffice-frontend.conf \
#           /etc/supervisor/conf.d/quantoffice-frontend.conf
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start quantoffice-frontend
#
# 端口：Vite dev 默认 5173（监听 0.0.0.0,已写死在 vite.config.ts）
#      如要和 nginx 共存,把 nginx 的 location / 改 proxy_pass 到 5173（见 README）
#
# 关键：直接跑 `bun run dev` 在老版 bun(< 1.1)上会报
#       "SyntaxError: Unexpected reserved word",因为 bun x 走 node CJS loader
#       加载 vite.js 的 top-level await。wrapper 脚本强制 `bun --bun run dev`。
#
# 互斥：
#   - 和 ./deploy.sh start（uvicorn）无冲突
#   - 和 quantoffice.conf（后端 supervisor）可同时跑
#   - 和 nginx 的 location / try_files 走法互斥（dev 模式要把 / 改 proxy）
# =============================================================================

[program:quantoffice-frontend]
; 用 wrapper 启动:自动找 bun + 强制 --bun runtime
; PROJECT_ROOT 决定 frontend/ 在哪里;如项目不在 /workspace,改成实际路径
command=/opt/quantoffice/scripts/start-frontend-dev.sh
directory=/workspace/frontend

; 运行用户（建议同后端用 www-data，但需保证该用户对 frontend/ 有读权限）
user=www-data

; 自启 & 崩溃自愈
autostart=true
autorestart=true
startretries=10
; Vite 启动较慢，给 15s
startsecs=15

; 整组杀（Vite 启动后会派生 esbuild / 预构建子进程，必须整组杀）
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

; 优先级（数字小优先）
priority=200

; 1H1G 机器友好：限制 Vite/Node 堆内存
; PROJECT_ROOT 改成本项目实际根目录(默认 /workspace,生产常见 /data/work/QuantOffice)
environment=
    PATH="/usr/local/bin:/usr/bin:/bin",
    PROJECT_ROOT="/workspace",
    NODE_OPTIONS="--max-old-space-size=512",
    CI="true"

; 日志
stdout_logfile=/var/log/supervisor/quantoffice-frontend.out.log
stderr_logfile=/var/log/supervisor/quantoffice-frontend.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=10

; HMR 会刷大量 stdout，标记 relevant=true 让 supervisorctl tail 能看到
stdout_logfile_relevant=true
