# deploy/supervisor/quantoffice-frontend.conf
# =============================================================================
# QuantOffice 前端 — Vite dev 服务器（supervisor 托管）
# -----------------------------------------------------------------------------
# 适用场景：
#   - 改了前端代码想即时热重载（Vite HMR）
#   - 1H1G 机器跑不动 bun run build（dev 模式内存友好）
#   - 不想每次改前端都跑 ./deploy.sh install
#
# 拷贝位置：
#   /etc/supervisor/conf.d/quantoffice-frontend.conf
#
# 启用：
#   sudo mkdir -p /var/log/supervisor
#   sudo cp deploy/supervisor/quantoffice-frontend.conf \
#           /etc/supervisor/conf.d/quantoffice-frontend.conf
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start quantoffice-frontend
#
# 端口：Vite dev 默认 5173（监听 0.0.0.0,已写死在 vite.config.ts）
#      如要和 nginx 共存,把 nginx 的 location / 改 proxy_pass 到 5173（见 README）
#
# 互斥：
#   - 和 ./deploy.sh start（uvicorn）无冲突
#   - 和 quantoffice.conf（后端 supervisor）可同时跑
#   - 和 nginx 的 location / try_files 走法互斥（dev 模式要把 / 改 proxy）
# =============================================================================

[program:quantoffice-frontend]
; 用 bun run dev 拉起 Vite（HMR 启用，--host 0.0.0.0 对外可访问）
; 注意：vite.config.ts 已固定 host=0.0.0.0 + port=5173，无需命令行覆盖
command=/root/.bun/bin/bun run dev
directory=/workspace/frontend

; 运行用户（建议同后端用 www-data，但需保证该用户对 frontend/ 有读权限）
user=www-data

; 自启 & 崩溃自愈
autostart=true
autorestart=true
startretries=10
; Vite 启动较慢，给 15s
startsecs=15

; 整组杀（Vite 启动后会派生 esbuild / 预构建子进程，必须整组杀）
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

; 优先级（数字小优先）
priority=200

; 1H1G 机器友好：限制 Vite/Node 堆内存
environment=
    PATH="/root/.bun/bin:/workspace/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    BUN_INSTALL="/root/.bun",
    NODE_OPTIONS="--max-old-space-size=512",
    CI="true"

; 日志
stdout_logfile=/var/log/supervisor/quantoffice-frontend.out.log
stderr_logfile=/var/log/supervisor/quantoffice-frontend.err.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=10

; HMR 会刷大量 stdout，标记 relevant=true 让 supervisorctl tail 能看到
stdout_logfile_relevant=true
